# config/backup/crontab
# Автоматические бэкапы для Bitrix
# Формат: минута час день месяц день_недели команда

# Бэкап базы данных каждый день в 2:00
0 2 * * * /bin/bash /scripts/backup.sh database

# Бэкап файлов каждый день в 3:00
0 3 * * * /bin/bash /scripts/backup.sh files

# Очистка старых бэкапов каждую неделю в воскресенье в 4:00
0 4 * * 0 /bin/bash /scripts/backup.sh cleanup

# Полный бэкап каждую субботу в 1:00
0 1 * * 6 /bin/bash /scripts/backup.sh full