# config/cron/crontab
# ============================================================================
# BITRIX CRONTAB CONFIGURATION
# Задачи для 1С-Битрикс: Управление сайтом
# ============================================================================
# Формат: минута час день месяц день_недели команда
# Логи автоматически попадают в /var/log/cron через supervisord

SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin
MAILTO=""

# ============================================================================
# БИТРИКС АГЕНТЫ - ОСНОВНАЯ ЗАДАЧА
# ============================================================================
# Запуск агентов каждую минуту
# Агенты выполняют отложенные операции: рассылки, индексацию, очистку и т.д.
* * * * * /usr/local/bin/php -f /home/bitrix/app/bitrix/modules/main/tools/cron_events.php

# ============================================================================
# ПОЧТОВАЯ ОЧЕРЕДЬ
# ============================================================================
# Отправка писем из очереди каждый час
# Битрикс складывает письма в очередь для оптимизации отправки
0 * * * * /usr/local/bin/php -f /home/bitrix/app/bitrix/modules/main/tools/mail_queue.php

# ============================================================================
# СОЦИАЛЬНЫЕ СЕТИ И ИНТЕГРАЦИИ
# ============================================================================
# Обработка очереди социальных сетей каждые 30 минут (если используется)
# */30 * * * * /usr/local/bin/php -f /home/bitrix/app/bitrix/modules/socialservices/classes/general/authmanager.php

# ============================================================================
# ПОИСКОВАЯ ИНДЕКСАЦИЯ
# ============================================================================
# Переиндексация поиска раз в день ночью (если не используется автоиндексация)
# 0 3 * * * /usr/local/bin/php -f /home/bitrix/app/bitrix/modules/search/bin/search_reindex.php

# ============================================================================
# ОЧИСТКА КЭША И ВРЕМЕННЫХ ФАЙЛОВ
# ============================================================================
# Очистка устаревшего кэша раз в день в 4 утра
0 4 * * * /usr/local/bin/php -f /home/bitrix/app/bitrix/modules/main/tools/cache_clean.php

# Очистка старых файлов в /tmp раз в неделю (воскресенье в 5 утра)
0 5 * * 0 find /home/bitrix/tmp -type f -mtime +7 -delete 2>/dev/null

# ============================================================================
# МОНИТОРИНГ И HEALTHCHECK
# ============================================================================
# Проверка работоспособности каждые 5 минут
# */5 * * * * curl -f http://localhost/healthcheck.php > /dev/null 2>&1 || echo "Site healthcheck failed"

# ============================================================================
# БЭКАПЫ (если нужны регулярные бэкапы БД на уровне приложения)
# ============================================================================
# Автобэкап БД каждый день в 2 ночи (раскомментируйте если нужно)
# 0 2 * * * /usr/local/bin/php -f /home/bitrix/app/local/scripts/backup_database.php

# ============================================================================
# LOGROTATE - РОТАЦИЯ ЛОГОВ
# ============================================================================
# Запуск logrotate каждый день в 23:55
# Логи хранятся 7 дней, затем удаляются
55 23 * * * /usr/sbin/logrotate /etc/logrotate.d/bitrix -s /var/lib/logrotate/status 2>&1 | logger -t logrotate

# ============================================================================
# КАСТОМНЫЕ ЗАДАЧИ
# ============================================================================
# Добавьте свои задачи ниже

# Пример: очистка логов старше 30 дней (каждый день в 6 утра)
# 0 6 * * * find /var/log/bitrix -name "*.log" -type f -mtime +30 -delete 2>/dev/null

# Пример: генерация отчетов (каждый понедельник в 9 утра)
# 0 9 * * 1 /usr/local/bin/php -f /home/bitrix/app/local/scripts/generate_reports.php