# config/cron/crontab
# ============================================================================
# ROOT CRONTAB CONFIGURATION FOR BITRIX CONTAINER
# NOTE: dcron in Alpine requires root crontab for reliable execution
# ============================================================================

# Alpine system maintenance tasks
# min	hour	day	month	weekday	command
*/15	*	*	*	*	run-parts /etc/periodic/15min
0	*	*	*	*	run-parts /etc/periodic/hourly
0	2	*	*	*	run-parts /etc/periodic/daily
0	3	*	*	6	run-parts /etc/periodic/weekly
0	5	1	*	*	run-parts /etc/periodic/monthly

# ============================================================================
# БИТРИКС АГЕНТЫ - ОСНОВНАЯ ЗАДАЧА
# ============================================================================
# Запуск агентов каждую минуту
# Агенты выполняют отложенные операции: рассылки, индексацию, очистку и т.д.
* * * * * /usr/local/bin/php -f /home/bitrix/app/bitrix/modules/main/tools/cron_events.php >> /var/log/cron/agents.log 2>&1

# ============================================================================
# ПОЧТОВАЯ ОЧЕРЕДЬ
# ============================================================================
# Отправка писем из очереди каждый час
0 * * * * /usr/local/bin/php -f /home/bitrix/app/bitrix/modules/main/tools/mail_queue.php >> /var/log/cron/mail_queue.log 2>&1

# ============================================================================
# ОЧИСТКА КЭША И ВРЕМЕННЫХ ФАЙЛОВ
# ============================================================================
# Очистка устаревшего кэша раз в день в 4 утра
0 4 * * * /usr/local/bin/php -f /home/bitrix/app/bitrix/modules/main/tools/cache_clean.php >> /var/log/cron/cache_clean.log 2>&1

# Очистка старых файлов в /tmp раз в неделю (воскресенье в 5 утра)
0 5 * * 0 find /home/bitrix/tmp -type f -mtime +7 -delete 2>/dev/null

# ============================================================================
# LOGROTATE - РОТАЦИЯ ЛОГОВ
# ============================================================================
# Запуск logrotate каждый день в 23:55
55 23 * * * /usr/sbin/logrotate /etc/logrotate.d/bitrix -s /var/lib/logrotate/status 2>&1 | logger -t logrotate

# ============================================================================
# КАСТОМНЫЕ ЗАДАЧИ
# ============================================================================
# Добавьте свои задачи ниже
