# config/cron/crontab
# ============================================================================
# ROOT CRONTAB CONFIGURATION FOR BITRIX CONTAINER
# NOTE: dcron in Alpine requires root crontab for reliable execution
# ============================================================================

# Alpine system maintenance tasks
# min	hour	day	month	weekday	command
*/15	*	*	*	*	run-parts /etc/periodic/15min
0	*	*	*	*	run-parts /etc/periodic/hourly
0	2	*	*	*	run-parts /etc/periodic/daily
0	3	*	*	6	run-parts /etc/periodic/weekly
0	5	1	*	*	run-parts /etc/periodic/monthly

# ============================================================================
# БИТРИКС АГЕНТЫ - МУЛЬТИСАЙТ
# ============================================================================
# Запуск агентов каждую минуту ДЛЯ ВСЕХ САЙТОВ
# Скрипт multisite-cron.sh обходит все домены и запускает агенты для каждого
* * * * * /usr/local/bin/scripts/multisite-cron.sh agents 2>&1 | logger -t bitrix-agents

# ============================================================================
# ПОЧТОВАЯ ОЧЕРЕДЬ - МУЛЬТИСАЙТ
# ============================================================================
# Отправка писем из очереди каждые 5 минут для всех сайтов
*/5 * * * * /usr/local/bin/scripts/multisite-cron.sh mail_queue 2>&1 | logger -t bitrix-mail

# ============================================================================
# ОЧИСТКА КЭША И ВРЕМЕННЫХ ФАЙЛОВ
# ============================================================================
# Очистка старых файлов в /tmp раз в неделю (воскресенье в 5 утра)
0 5 * * 0 find /home/bitrix/tmp -type f -mtime +7 -delete 2>/dev/null

# Ротация логов мультисайтового cron
0 0 * * * find /var/log/cron -name "*.log" -mtime +7 -delete 2>/dev/null

# ============================================================================
# LOGROTATE - РОТАЦИЯ ЛОГОВ
# ============================================================================
# Запуск logrotate каждый день в 23:55
55 23 * * * /usr/sbin/logrotate /etc/logrotate.d/bitrix -s /var/lib/logrotate/status 2>&1 | logger -t logrotate

# ============================================================================
# КАСТОМНЫЕ ЗАДАЧИ
# ============================================================================
# Добавьте свои задачи ниже
