# docker/php/bitrix/supervisord.conf
# ============================================================================
# SUPERVISOR CONFIGURATION FOR BITRIX CONTAINER
# Управление всеми процессами в едином Bitrix-контейнере
# ============================================================================

[unix_http_server]
file=/var/run/supervisor/supervisor.sock
chmod=0700
chown=%(ENV_UGN)s:%(ENV_UGN)s

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisor/supervisord.pid
childlogdir=/var/log/supervisor
loglevel=info
# Перенаправление логов в stdout/stderr для Promtail
logfile_maxbytes=0

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor/supervisor.sock

# ============================================================================
# PHP-FPM - ОСНОВНОЙ ПРОЦЕСС ОБРАБОТКИ PHP
# Master process runs as root (required for error_log to /dev/stderr)
# Worker processes run as bitrix (configured in www.conf)
# ============================================================================
[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize --fpm-config /usr/local/etc/php-fpm.conf
user=root
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
priority=1
startsecs=10
stopwaitsecs=60
# Graceful stop
stopsignal=QUIT
# Health check каждые 60 секунд
numprocs=1

# ============================================================================
# CRON - ПЛАНИРОВЩИК ЗАДАЧ ДЛЯ БИТРИКС АГЕНТОВ
# ============================================================================
[program:crond]
command=/usr/sbin/crond -f -l 2 -L /dev/stdout
user=root
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
priority=2
startsecs=5
stopwaitsecs=10
stopsignal=TERM

# ============================================================================
# BITRIX PUSH-SERVER (опционально, если используется)
# ============================================================================
# Раскомментируйте, если используете Push & Pull сервер Битрикса
# [program:bitrix-push-server]
# command=/usr/local/bin/php /home/%(ENV_UGN)s/app/bitrix/modules/pull/install/push-server/push-server.php
# user=%(ENV_UGN)s
# stdout_logfile=/var/log/supervisor/bitrix-push-server.log
# stderr_logfile=/var/log/supervisor/bitrix-push-server.err.log
# autostart=false
# autorestart=true
# priority=10
# startsecs=10
# stopwaitsecs=30

# ============================================================================
# BITRIX AGENTS (фоновые задачи, опционально)
# ============================================================================
# Альтернативный способ запуска агентов (если не используете cron)
# [program:bitrix-agents]
# command=/usr/local/bin/php /home/%(ENV_UGN)s/app/bitrix/modules/main/tools/cron_events.php
# user=%(ENV_UGN)s
# stdout_logfile=/var/log/supervisor/bitrix-agents.log
# stderr_logfile=/var/log/supervisor/bitrix-agents.err.log
# autostart=false
# autorestart=true
# priority=11
# startsecs=5
# stopwaitsecs=60

# ============================================================================
# BITRIX MAIL QUEUE (очередь отправки почты, опционально)
# ============================================================================
# [program:bitrix-mail-queue]
# command=/usr/local/bin/php /home/%(ENV_UGN)s/app/bitrix/modules/main/tools/mail_queue.php
# user=%(ENV_UGN)s
# stdout_logfile=/var/log/supervisor/bitrix-mail-queue.log
# stderr_logfile=/var/log/supervisor/bitrix-mail-queue.err.log
# autostart=false
# autorestart=true
# priority=12
# startsecs=5
# stopwaitsecs=30

# ============================================================================
# ВКЛЮЧЕНИЕ ДОПОЛНИТЕЛЬНЫХ КОНФИГУРАЦИЙ
# ============================================================================
# Загружаем пользовательские программы из /etc/supervisor/conf.d
[include]
files = /etc/supervisor/conf.d/*.conf
