# scripts/bitrix-docker.service
# ============================================================================
# SYSTEMD SERVICE FOR BITRIX DOCKER
# Auto-start containers after server reboot with health verification
#
# Installation:
#   sudo ./scripts/install-service.sh install
#
# Or use:
#   make install-service
# ============================================================================

[Unit]
Description=Bitrix Docker Environment
Documentation=https://github.com/your-repo/b-docker
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service
# Prevent rapid restart loops: max 3 attempts in 60 seconds
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
Group=docker

# Working directory (replaced by install script)
WorkingDirectory=/home/deploy/bitrix

# Environment file for COMPOSE_PROFILES, DOMAIN, etc.
EnvironmentFile=-/home/deploy/bitrix/.env

# Wait for Docker daemon to be fully ready before starting
ExecStartPre=/bin/bash -c 'until docker info >/dev/null 2>&1; do sleep 2; done'

# Create Docker network if it does not exist yet
ExecStartPre=/bin/bash -c 'docker network inspect ${COMPOSE_PROJECT_NAME:-bitrix}_network >/dev/null 2>&1 || docker network create ${COMPOSE_PROJECT_NAME:-bitrix}_network --subnet=172.20.0.0/16 || true'

# Start containers with all production profiles
# NOTE: Profiles can be overridden via COMPOSE_PROFILES in .env
ExecStart=/usr/bin/docker compose -f docker-compose.bitrix.yml --profile prod --profile monitoring --profile backup --profile push --profile security up -d --remove-orphans

# Health verification: wait for containers, then check health status
ExecStartPost=/bin/bash -c 'sleep 10 && docker compose -f docker-compose.bitrix.yml ps --format "{{.Name}} {{.Status}}" | grep -q "healthy" && echo "Services are healthy" || echo "WARNING: Some services may not be healthy yet"'

# Stop containers gracefully with extended timeout for DB flush
ExecStop=/usr/bin/docker compose -f docker-compose.bitrix.yml --profile prod --profile monitoring --profile backup --profile push --profile security down --timeout 60

# Restart on failure only (not on clean stop)
Restart=on-failure
RestartSec=30

# Generous timeouts for large deployments with DB warmup
TimeoutStartSec=600
TimeoutStopSec=180

# Logging to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bitrix-docker

[Install]
WantedBy=multi-user.target
